# spline-core.service
[Unit]
Description=Spline Core — AI Training Prescription Daemon
Documentation=file:///DESIGN.md
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
ExecStart=/usr/local/bin/spline
WorkingDirectory=/var/lib/spline
EnvironmentFile=/etc/spline/spline.env

# Static user configuration for external database access
User=spline
Group=spline
StateDirectory=spline
StateDirectoryMode=0755
ConfigurationDirectory=spline

# Crash recovery (DESIGN §8)
Restart=on-failure
RestartSec=30

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spline-core

[Install]
WantedBy=multi-user.target
